# packages/sysutils/systemd/system.d/serial-getty@.service.d/auto-login.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -i --skip-login --login-program /bin/bash -o '-l' --noclear --keep-baud  %I $TERM
